# feh/update-desktop-background.timer
[Unit]
Description=automatically changing desktop background

[Timer]
OnActiveSec=1h
OnBootSec=5s

[Install]
WantedBy=timers.target
